# data/mq_unipi_daemon
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/local/bin/mq_unipi_daemon"
command_args=""
pidfile="/var/run/mq_unipi_daemon.pid"

depend() {
 use net pigpiod mosquitto
}
